# deploy/polymarket-forward-monitor.timer.example
[Unit]
Description=Run Polymarket forward-test monitor every five minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
