# deploy/polymarket-backup.timer.example
[Unit]
Description=Create daily Polymarket production backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
